# opengeek-llm-chat.service
[Unit]
Description=OpenGeekLLMChat Server (llama.cpp backend)
After=network.target
Wants=network.target

[Service]
Type=simple

# ⚠️ 以下を環境に合わせて変更してください
User=wizapply-ai
Group=wizapply-ai
WorkingDirectory=/home/wizapply-ai/opengeek-llm-chat
ExecStart=/usr/bin/node /home/wizapply-ai/opengeek-llm-chat/server.js

# 再起動設定
Restart=always
RestartSec=10

# 環境変数（必要に応じて）
Environment=NODE_ENV=production
Environment=PORT=3000

# AMD ROCm環境では HSA_OVERRIDE_GFX_VERSION 等が必要な場合あり
# Environment=HSA_OVERRIDE_GFX_VERSION=12.0.1
# Environment=HIP_VISIBLE_DEVICES=0,1

# ログ出力
StandardOutput=journal
StandardError=journal

# セキュリティ
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

# 書き込み許可ディレクトリ
# - /home/wizapply-ai/opengeek-llm-chat: chats/ public/uploads/ public/plots/ settings.json 等
# - /tmp: Pythonスクリプト一時ファイル
ReadWritePaths=/home/wizapply-ai/opengeek-llm-chat /tmp

# llama-server を spawn するためバイナリ実行が必要
# Pythonスクリプト実行も含む
PrivateTmp=false

[Install]
WantedBy=multi-user.target
